# tools/ansible/playbook/roles/setup_env/templates/couscous_app.service.j2
[Unit]
Description=Couscous app service

[Service]
User=pi
ExecStart=/home/{{ PI_USERNAME }}/application_test_gui_exe
Restart=always