# ansible/playbook/roles/setup_pi/files/couscous_app.service
[Unit]
Description=Couscous app service

[Service]
User=pi
ExecStart=/home/pi/application_test_gui_exe
Restart=always